# docker/local/s6-rc.d/postgres/finish
#!/command/with-contenv bash
exec /usr/local/bin/specgate-service-finish postgres "$@"
